# infra/launchd/com.medinovai.ceo-stack.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchDaemon: com.medinovai.ceo-stack
  ─────────────────────────────────────────────────────────────────────────────
  Starts the authoritative AtlasOS prod environment and the Dockerized CEO
  gateway stack on every boot, BEFORE any user logs in. Restarts automatically
  on crash.

  Install (run once as root):
    sudo cp /Users/mayanktrivedi/Github/medinovai-health/medinovai-Deploy/infra/launchd/com.medinovai.ceo-stack.plist \
         /Library/LaunchDaemons/com.medinovai.ceo-stack.plist
    sudo launchctl load -w /Library/LaunchDaemons/com.medinovai.ceo-stack.plist

  Check status:
    sudo launchctl list | grep medinovai
    tail -f /var/log/medinovai-ceo-stack.log

  Uninstall:
    sudo launchctl unload /Library/LaunchDaemons/com.medinovai.ceo-stack.plist
    sudo rm /Library/LaunchDaemons/com.medinovai.ceo-stack.plist
  ─────────────────────────────────────────────────────────────────────────────
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.medinovai.ceo-stack</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/mayanktrivedi/Github/medinovai-health/medinovai-Deploy/infra/launchd/start-ceo-stack.sh</string>
  </array>

  <!-- Start at boot, before login -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Wait for Docker to be available (Docker Desktop starts its socket after ~30s) -->
  <key>StartInterval</key>
  <integer>0</integer>

  <!-- Restart if the process exits unexpectedly -->
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <!-- Wait 30 seconds before restarting on failure (gives Docker time to start) -->
  <key>ThrottleInterval</key>
  <integer>30</integer>

  <!-- Log output -->
  <key>StandardOutPath</key>
  <string>/var/log/medinovai-ceo-stack.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/medinovai-ceo-stack.log</string>

  <!-- Environment: PATH must include docker -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Applications/Docker.app/Contents/Resources/bin</string>
    <key>HOME</key>
    <string>/Users/mayanktrivedi</string>
    <key>ATLASOS_PATH</key>
    <string>/Users/mayanktrivedi/Github/medinovai-health/AtlasOS</string>
    <key>ATLAS_COMMAND_PATH</key>
    <string>/Users/mayanktrivedi/Github/_non-org/atlas-command</string>
  </dict>

  <!-- Working directory -->
  <key>WorkingDirectory</key>
  <string>/Users/mayanktrivedi/Github/medinovai-health/medinovai-Deploy</string>
</dict>
</plist>
